# install/gpio-server.conf
description "GPIO Server"
start on runlevel [2345]
stop on runlevel [06]
respawn
respawn limit 10 5

script
        cd /home/pi/Pi-GPIO-Server
        sudo python pi_gpio_server.py
end script
